# service/update-unbound-hosts.service
[Unit]
Description=Reboot

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hosts2unbound -o /etc/unbound/auto-block-zones.conf
